# src/program-statsd.conf
; -*- ini -*-
[program:influxdb]
user = nobody
command = /usr/local/bin/statsd /etc/statsd_config.js

redirect_stderr = true

stdout_logfile = /var/log/statsd.log
stdout_logfile_maxbytes = 5MB
